# modules/home/uair/uair.service
[Unit]
Description=Uair

[Service]
Type=simple
ExecStart=/usr/bin/env uair -q

[Install]
WantedBy=default.target
